# services/wifi-watchdog.service
[Unit]
Description=WiFi Watchdog Service
After=network.target

[Service]
ExecStart=/usr/local/bin/wifi-watchdog.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

